# sdk_aa/target/linux/ar71xx/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
ttyS0::askfirst:/bin/login

# Guarantee that RFMGMT always execute
#::respawn:/usr/bin/rfmgmt
